# docker/common/etc/crontab
#
# Regenerating the optimized pages for SEO
#
00 */8 * * * cd /www/frontend-prerenderer && make regenerate_cached_pages
